# thinkfan/etc/systemd/system/thinkfan.service.d/override.conf
[Service]
Environment='THINKFAN_ARGS=-b0'
