# systemd/.config/systemd/user/dotfiles-maintenance.service
[Unit]
Description=Dotfiles user cache maintenance

[Service]
Type=oneshot
ExecStart=%h/.local/bin/system-maintenance user-cache
