# displaylink/etc/systemd/system/displaylink.service.d/override.conf
[Service]
ExecStartPre=
ExecStartPre=/bin/sh -c 'grep -q "^evdi " /proc/modules || /sbin/modprobe evdi'
